# etc/supervisor/conf.d/supervisor.conf
[supervisord]
nodaemon=true

[program:zabbix-agent]
command=/usr/sbin/zabbix_agentd -f
