# etc/zabbix/crontab
*/1 * * * * /bin/sh /etc/zabbix/scripts/trapper/docker
